# crontab
@reboot /usr/bin/python /home/pi/Desktop/autoRunCollection.py